# profile.d/aliases.sh
#!/usr/bin/env bash
# Aliases
alias cat="bat"
alias tig="tig --show-signature"
alias iso_date="date "+%Y-%m-%d""
